# devices/additions/etc/init.d/adbd
#!/sbin/openrc-run
name="Android Debugging Bridge Daemon"
description="Debugging utility."

command="/sbin/adbd"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need localmount
	after bootmisc
}

start_pre() {
	mountpoint /sys/kernel/config/ || mount -t configfs configfs /sys/kernel/config
}